# hadoop/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:thrift]
command=/usr/bin/bash -c "/root/run_hadoop.sh"