# zeppelin/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:zeppl]
command=/usr/bin/bash -c "/root/run_zeppl.sh"