# kafka/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:kafka]
command=/bin/bash -c "/root/run_kafka.sh"